# states/flatpak/files/flatpak-cleanup.service
[Unit]
Description=Cleanup flatpak packets ({{ flatpak_mode }})

[Service]
ExecStart=/usr/bin/flatpak --{{ flatpak_mode }} uninstall --unused
